# .docker/config/nginx.supervisord.conf
[supervisord]
nodaemon=true

[program:sysctl]
command=sysctl -p

[program:keepalived]
command=service keepalived start

[program:nginx]
command=nginx -g 'daemon off;'